# data/tcmmd.service
[Unit]
Description=Traffic control multimedia daemon

[Service]
ExecStart=/usr/bin/tcmmd
BusName=org.tcmmd

[Install]
Alias=dbus-org.tcmmd.service
